# roles/web/templates/supervisor_ngrok.conf.j2
[program:ngrok]
command=/usr/local/bin/ngrok http --log stdout {{ app_port }}
user=root
stdout_logfile=/var/log/supervisor/ngrok_moonlight.log
redirect_stderr=true
